# etc/periodic/daily/looper
#!/bin/sh

set -e

set -o allexport
source /cron-env
set +o allexport

cd /opt/app
exec s6-setuidgid "${PUID}:${PGID}" ./looper.sh
