# etc/periodic/daily/looper
#!/bin/sh

source /cron-env

cd /opt/app
s6-setuidgid '${PUID}:${PGID}' ./looper.sh
