# etc/periodic/daily/looper
#!/bin/sh

set -e

source /cron-env

cd /opt/app
exec s6-setuidgid "${PUID}:${PGID}" ./looper.sh
