# etc/s6/backend/run
#!/bin/sh

set -e

cd /opt/app
# TODO: Make the API
#exec s6-setuidgid "${PUID}:${PGID}" node ./dist/api/index.js
sleep 3600
