# etc/s6/cron/run
#!/bin/sh

set -e

touch /cron-env
chmod 400 /cron-env
printenv > /cron-env

exec crond -s -n
